# scripts/zhatroom.service
[Unit]
Description=ZhatRoom Chat Server
After=network.target docker.service

[Service]
Type=simple
User=chat
Group=chat
WorkingDirectory=/opt/zhatroom
ExecStart=/opt/zhatroom/bin/server
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
