# scripts/zhatroom.service
[Unit]
Description=ZhatRoom Chat Server
After=network.target docker.service

[Service]
Type=simple
User=zhat
Group=zhat
WorkingDirectory=/opt/zhatroom
ExecStart=/opt/zhatroom/bin/server
Restart=always
RestartSec=5

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/zhatroom /tmp

[Install]
WantedBy=multi-user.target
